# systemd/codex-connect-gateway.service.template
[Unit]
Description=Codex Connect Gateway
After=codex-connect-app-server.service
Wants=codex-connect-app-server.service

[Service]
Type=simple
WorkingDirectory=__CONFIG_DIR__
ExecStart=__NODE_BINARY__ __CLI_ENTRY__ gateway
Environment="CODEX_CONNECT_HOME=__CONFIG_DIR_ENV__"
Environment="CODEX_CONNECT_CONFIG_FILE=__CONFIG_PATH_ENV__"
Environment="CODEX_BINARY=__CODEX_BINARY_ENV__"
Environment="PATH=__SYSTEMD_PATH__"
UMask=0077
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
