# systemd/codex-connect-app-server.service.template
[Unit]
Description=Codex Connect App Server

[Service]
Type=simple
WorkingDirectory=__WORKDIR__
ExecStart=__NODE_BINARY__ __CLI_ENTRY__ service-app-server
Environment="CODEX_CONNECT_HOME=__CONFIG_DIR__"
Environment="CODEX_CONNECT_CONFIG_FILE=__CONFIG_PATH_ENV__"
Environment="PATH=__SYSTEMD_PATH__"
UMask=0077
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
