# launchd/com.msola.codex-gateway.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.msola.codex-gateway</string>
  <key>ProgramArguments</key>
  <array>
    <string>__NODE_BINARY__</string>
    <string>__CLI_ENTRY__</string>
    <string>gateway</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__CONFIG_DIR__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>CODEX_CONNECT_HOME</key>
    <string>__CONFIG_DIR__</string>
    <key>CODEX_CONNECT_ENV_FILE</key>
    <string>__ENV_PATH__</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>5</integer>
  <key>StandardOutPath</key>
  <string>__RUNTIME_DIR__/gateway.log</string>
  <key>StandardErrorPath</key>
  <string>__RUNTIME_DIR__/gateway.error.log</string>
</dict>
</plist>
