# launchd/com.hegenai.codex-app-server.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hegenai.codex-app-server</string>
  <key>ProgramArguments</key>
  <array>
    <string>__CODEX_BINARY__</string>
    <string>app-server</string>
    <string>--listen</string>
    <string>unix://__SOCKET_PATH__</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__WORKDIR__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__LAUNCHD_PATH__</string>
    <key>HTTP_PROXY</key>
    <string>__HTTP_PROXY__</string>
    <key>HTTPS_PROXY</key>
    <string>__HTTPS_PROXY__</string>
    <key>ALL_PROXY</key>
    <string>__ALL_PROXY__</string>
    <key>NO_PROXY</key>
    <string>__NO_PROXY__</string>
    <key>http_proxy</key>
    <string>__http_proxy__</string>
    <key>https_proxy</key>
    <string>__https_proxy__</string>
    <key>all_proxy</key>
    <string>__all_proxy__</string>
    <key>no_proxy</key>
    <string>__no_proxy__</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>5</integer>
  <key>StandardOutPath</key>
  <string>__RUNTIME_DIR__/codex-app-server.log</string>
  <key>StandardErrorPath</key>
  <string>__RUNTIME_DIR__/codex-app-server.error.log</string>
</dict>
</plist>
